# super.conf
[unix_http_server]
file=/tmp/supervisor-convert.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-convert.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-convert.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-convert.sock ; use a unix:// URL  for a unix socket

[program:gunicorn-convert]
command=newrelic-admin run-program gunicorn -c gunicorn.conf.py convert.wsgi
environment=NEW_RELIC_CONFIG_FILE="../newrelic.ini"
autostart=true                ; start at csupervisord start (default: true)
stdout_logfile=../gunicorn-out.log        ; stdout log path, none for none; default auto
stderr_logfile=../gunicorn-err.log        ; stderr log path, NONE for none; default AUTO

[program:celeryd]
command=celery worker --app=convert --loglevel=info
directory=./convert
autostart=true                ; start at csupervisord start (default: true)
stdout_logfile=../celery-out.log        ; stdout log path, none for none; default auto
stderr_logfile=../celery-err.log        ; stderr log path, NONE for none; default AUTO